# docker/dockerfiles/base/rootfs/etc/profile.d/010-jvm.sh
export JAVA_HOME=/usr/local/jvm/java-8-openjdk-amd64
PATH=$PATH:$JAVA_HOME/bin
